# templates/consul-template.upstart.erb
# Consul Agent (Upstart unit)
description "Consul Agent"
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

env CONSUL=<%= scope.lookupvar('consul_template::bin_dir') %>/consul-template
env CONFIG=<%= scope.lookupvar('consul_template::config_dir') %>


script
    exec $CONSUL -config $CONFIG/config <%= scope.lookupvar('consul_template::extra_options') %>
end script

respawn
respawn limit 10 10
kill timeout 10
